# etc/rc.local
printf "nameserver 172.16.12.6\nnameserver 172.16.12.5\nnameserver 172.16.12.1\n" >> /etc/resolv.conf
route delete default
route add default 172.16.12.6
