# etc/cron.d/mount
* * * * *	root	/sbin/mount | /usr/bin/grep -q 's3fs' || /sbin/mount -a
